# loadout/cheats/11-system_and_network_tools/autorun-organizer.yaml
tool: autorun-organizer
tags:
- m11
- network
- cli
exam_hint: 'Autorun Organizer is a Windows startup management utility that displays and manages all programs configured to run automatically when Windows starts. It covers registry Run keys, startup folders, scheduled tasks, services, browser extensions, and other autostart locations. Security analysts use autorun management tools to detect malware persistence mechanisms, as virtually all Windows malware esta Strategy: nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on th'
docs:
- https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns
- https://attack.mitre.org/techniques/T1547/
actions:
- title: Autorunsc.exe -a *
  desc: Runs `autorunsc.exe` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: autorunsc.exe -a * -user * -h -c > autoruns_baseline.csv
- title: Get-ItemProperty HKLM:\SOFTWARE\Microsoft\Windows\CurrentVer
  desc: Runs `Get-ItemProperty` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: Get-ItemProperty HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run
- title: Get-ItemProperty HKCU:\SOFTWARE\Microsoft\Windows\CurrentVer
  desc: Runs `Get-ItemProperty` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: Get-ItemProperty HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run
- title: Get-ItemProperty HKLM:\SOFTWARE\Microsoft\Windows\CurrentVer (2)
  desc: Runs `Get-ItemProperty` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: Get-ItemProperty HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce
- title: Dir "C:\ProgramData\Microsoft\Windows\Start Menu\Programs\St
  desc: Runs `dir` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: dir "C:\ProgramData\Microsoft\Windows\Start Menu\Programs\Startup"
- title: Dir "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Star
  desc: Runs `dir` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: dir "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Startup"
