# loadout/cheats/11-system_and_network_tools/started-pro.yaml
tool: started-pro
tags:
- m11
- network
- cli
exam_hint: 'StartEd Pro is a Windows startup management tool that provides comprehensive visibility and management of all programs, scripts, and processes that automatically execute at Windows startup. It covers registry autorun keys (HKLM/HKCU Run, RunOnce), startup folders, scheduled tasks, Win.ini, System.ini, and other autostart locations. Analysts use startup management tools like StartEd Pro to detect a Strategy: nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on th'
docs:
- https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns
- https://attack.mitre.org/techniques/T1547/
actions:
- title: Autorunsc.exe -accepteula -a
  desc: Runs `autorunsc.exe` as part of the CEH workflow. nc/socat for relays; curl/wget for payload staging; tcpdump/tshark with BPF filters on the pivot interface to limit noise. Authorized lab targets only.
  command: autorunsc.exe -accepteula -a * -h -c > autoruns_output.csv
